# packaging/cbhm.path
[Unit]
Description=Start the "cbhm" program

[Path]
PathExists=/run/.wm_ready

