# packaging/cbhm.service
[Unit]
Description=Clipboard History Manager

[Service]
Type=simple
ExecStart=/usr/apps/org.tizen.cbhm/bin/cbhm
